# lenovod330-10igl-webcam.service
[Unit]
Description=Change default webcam

[Service]
ExecStart=/usr/bin/lenovod330-10igl-webcam.sh

[Install]
WantedBy=multi-user.target
